# systemd-automount-example/media-share01.automount
[Unit]
Description=Share 01

[Automount]
Where=/media/share01

[Install]
WantedBy=multi-user.target
